# deploy/ec2/systemd/two-do-auth.service.example
# Install: sudo cp two-do-auth.service.example /etc/systemd/system/two-do-auth.service
# Edit WorkingDirectory, ExecStart, User=, and EnvironmentFile to match your server, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now two-do-auth
#
# User: Amazon Linux 2023 — ec2-user. Ubuntu — often www-data (chown secrets accordingly).
# Node: use /usr/bin/node (dnf) or full path to nvm’s node if you do not install system Node.

[Unit]
Description=two-do Firebase session (auth-server)
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/two-do/services/auth-server
Environment=NODE_ENV=production
EnvironmentFile=/opt/two-do/services/auth-server/.env
ExecStart=/usr/bin/node /opt/two-do/services/auth-server/src/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

# Nested clone (repo at /opt/two-do/two-do): set WorkingDirectory, EnvironmentFile, and
# ExecStart to paths under .../two-do/services/auth-server/ instead.
