# deploy/ec2/systemd/two-do-auth.service.example
# Install: sudo cp two-do-auth.service.example /etc/systemd/system/two-do-auth.service
# Edit paths and User=, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now two-do-auth

[Unit]
Description=two-do Firebase session (auth-server)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/two-do/services/auth-server
Environment=NODE_ENV=production
EnvironmentFile=/opt/two-do/services/auth-server/.env
ExecStart=/usr/bin/node /opt/two-do/services/auth-server/src/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
